# etc/systemd/system/checkpoint401.service
[Unit]
Description=checkpoint401
After=network.target

[Service]
ExecStart=/usr/bin/deno run --allow-net=127.0.0.1:3000 --allow-read=/opt/checkpoint401 --allow-env=PORT,LISTEN_ADDRESS --allow-write=/opt/checkpoint401/config/route_stats_counters.db /opt/checkpoint401/checkpoint401.ts
WorkingDirectory=/opt/checkpoint401/config
Restart=always
User=ubuntu
Environment=PATH=/usr/bin:/usr/local/bin

# Hardening. Adjust ReadWritePaths to point at your --db-filename.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=no
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
ReadWritePaths=/opt/checkpoint401/config

[Install]
WantedBy=multi-user.target
